# debian/postinst
#!/bin/sh

set -e

. /etc/default/openmediavault
. /usr/share/openmediavault/scripts/helper-functions

case "$1" in
    configure)
        # Set the default settings of the service package to those expected by
        # OpenMediaVault.
        if [ -z "${2}" ]; then
            deb-systemd-helper disable mysql.service >/dev/null || true
            deb-systemd-invoke stop mysql.service >/dev/null || true
        fi

        # Configuration.
        SERVICE_XPATH_NAME="mysql"
        SERVICE_XPATH="/config/services/${SERVICE_XPATH_NAME}"

        if ! omv_config_exists "${SERVICE_XPATH}"; then
            omv_config_add_element "/config/services" "${SERVICE_XPATH_NAME}"
            omv_config_add_element "${SERVICE_XPATH}" "enable" "0"
            omv_config_add_element "${SERVICE_XPATH}" "enable_networking" "1"
            omv_config_add_element "${SERVICE_XPATH}" "port" "3306"
            omv_config_add_element "${SERVICE_XPATH}" "bind_address" "127.0.0.1"
            omv_config_add_element "${SERVICE_XPATH}" "disable_aio" "0"
            omv_config_add_element "${SERVICE_XPATH}" "extra_options" ""
            omv_config_add_element "${SERVICE_XPATH}" "enable_management_site" "0"
        fi

        # Activate package triggers. These triggers are only set during the
        # package installation.
        dpkg-trigger update-fixperms
        dpkg-trigger update-locale
    ;;

    abort-upgrade|abort-remove|abort-deconfigure)
    ;;

    *)
        echo "postinst called with unknown argument" >&2
        exit 1
    ;;
esac

#DEBHELPER#

exit 0
